# docker/supervisorConfigs/widgets.conf
[program:widgets]
directory=/usr/src/app/widgets
command=grunt connect
autostart=true
autorestart=true
stderr_logfile=/var/log/widgets.err.log
stdout_logfile=/var/log/widgets.out.log